# com.btcomega.trader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.btcomega.trader</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>cd /Users/bollineni/btc-omega-vi &amp;&amp; ./run_trader_cycle.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/bollineni/btc-omega-vi</string>
	<key>StartInterval</key>
	<integer>300</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/bollineni/btc-omega-vi/trader.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/bollineni/btc-omega-vi/trader.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
</dict>
</plist>
